# launchd/com.vincent.idea-wheel-keepalive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vincent.idea-wheel-keepalive</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/vincent/.openclaw/workspace/projects/idea-wheel/scripts/keep-awake.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>WorkingDirectory</key>
  <string>/Users/vincent/.openclaw/workspace/projects/idea-wheel</string>

  <key>StandardOutPath</key>
  <string>/Users/vincent/.openclaw/workspace/projects/idea-wheel/logs/keep-alive.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/vincent/.openclaw/workspace/projects/idea-wheel/logs/keep-alive.error.log</string>
</dict>
</plist>
